# Auto Install/config/atlas-ai-healthcheck.timer
[Unit]
Description=Atlas AI Support Assistant Health Check Timer
Documentation=https://github.com/your-username/atlas-ai
Requires=atlas-ai-healthcheck.service

[Timer]
# Run every 5 minutes
OnCalendar=*:0/5
Persistent=true

# Start immediately after boot
OnBootSec=5min

[Install]
WantedBy=timers.target