# Auto Install/config/atlas-ai-backup.service
[Unit]
Description=Atlas AI Support Assistant Backup Service
Documentation=https://github.com/your-username/atlas-ai

[Service]
Type=oneshot
User=atlas
Group=atlas
WorkingDirectory=/etc/atlas-ai

# Backup command
ExecStart=/etc/atlas-ai/backup.sh

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/atlas-ai /var/log/atlas-ai
CapabilityBoundingSet=

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-ai-backup