# auto-install/config/atlas-ai-backup.timer
[Unit]
Description=Atlas AI Support Assistant Backup Timer
Documentation=https://github.com/leipnar/Atlas-AI
Requires=atlas-ai-backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=*-*-* 02:00:00
Persistent=true

# Randomize start time to avoid load spikes
RandomizedDelaySec=300

[Install]
WantedBy=timers.target