# auto-install/config/atlas-ai-healthcheck.service
[Unit]
Description=Atlas AI Support Assistant Health Check Service
Documentation=https://github.com/your-username/atlas-ai

[Service]
Type=oneshot
User=atlas
Group=atlas
WorkingDirectory=/etc/atlas-ai

# Health check command
ExecStart=/etc/atlas-ai/health-check.sh

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/atlas-ai
CapabilityBoundingSet=

# Resource limits
MemoryLimit=128M
CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-ai-healthcheck